# config/cherrypy-server.conf
description "CherryPy webserver"
start on vagrant-mounted

script
    chdir /vagrant/demo
    /usr/bin/python server.py
end script
